# spec/integration/vm/cookbooks/app/files/default/rackup.conf
[program:rackup]
directory=/vagrant/app
command=rackup -p 80 config.ru
