# spec/integration/vm/recipes/app/files/etc/supervisor/conf.d/rackup.conf
[program:rackup]
directory=/vagrant/app
command=rackup -o 0.0.0.0 -p 80 config.ru
